# src/record_date.service
[Unit]
Description=Record data in a file

[Service]
ExecStart=/tmp/record_date.sh